# archiso/profiles/playos/airootfs/etc/systemd/system/playos-firstboot.service
[Unit]
Description=Pre-configure system settings on first boot
Before=systemd-firstboot.service
ConditionFirstBoot=yes

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/timedatectl set-timezone UTC
ExecStart=/usr/bin/localectl set-locale LANG=en_US.UTF-8
ExecStart=/usr/bin/localectl set-keymap us
ExecStart=/usr/bin/systemctl enable sshd
ExecStart=/usr/bin/systemctl start sshd

[Install]
WantedBy=sysinit.target
