# archiso/profiles/playos/airootfs/etc/systemd/system/playos-erofs-fixup.service
[Unit]
Description=Fix EROFS-corrupted firmware and modules for GPU
DefaultDependencies=no
Before=systemd-udev-trigger.service
Before=playos-compositor.service
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/playos/erofs-fixup.sh

[Install]
WantedBy=sysinit.target
