# archiso/profiles/playos/airootfs/etc/systemd/system/playos-game.slice
[Unit]
Description=PlayOS Active Game Slice

[Slice]
CPUWeight=8000
IOWeight=8000
MemoryHigh=85%
