# archiso/profiles/playos/airootfs/etc/systemd/system/playos-library.service
[Unit]
Description=PlayOS Game Library Scanner
After=playos-compositor.service
PartOf=playos-async.target

[Service]
Type=simple
User=playos
Group=playos
ExecStart=/usr/bin/bash -lc 'echo "[PlayOS] library scan placeholder"; sleep infinity'
Restart=on-failure
Nice=10
CPUWeight=500
IOWeight=500

[Install]
WantedBy=playos-async.target
