# archiso/profiles/playos/airootfs/etc/systemd/system/playos-ui.slice
[Unit]
Description=PlayOS UI and Compositor Slice

[Slice]
CPUWeight=10000
IOWeight=10000
MemoryHigh=80%
