# archiso/profiles/playos/airootfs/etc/systemd/system/playos-compositor.service
[Unit]
Description=PlayOS Compositor
After=local-fs.target systemd-udevd.service systemd-udev-trigger.service seatd.service
Wants=seatd.service
Before=playos-async.target

[Service]
Type=simple
User=playos
Group=playos
WorkingDirectory=/var/lib/playos

Environment=XDG_RUNTIME_DIR=/run/playos
Environment=PLAYOS_SHELL_CMD=/usr/bin/playos-shell

# GPU-adaptive renderer: pixman for virtual GPUs, EGL/GLES2 for real hardware.
# The ExecStartPre script writes /run/playos/compositor.env at boot.
ExecStartPre=/usr/bin/playos-compositor-env
EnvironmentFile=-/run/playos/compositor.env

RuntimeDirectory=playos
RuntimeDirectoryMode=0755

ExecStart=/usr/bin/playos-compositor /usr/bin/playos-shell

Restart=always
RestartSec=1

Slice=playos-ui.slice
Nice=-10
CPUWeight=10000
IOWeight=10000
OOMScoreAdjust=-900
LimitMEMLOCK=infinity
LimitRTPRIO=20

[Install]
WantedBy=playos-visual.target
