# archiso/profiles/playos/airootfs/etc/systemd/system/playos-bluetooth.service
[Unit]
Description=PlayOS Bluetooth (BlueZ)
# Start AFTER compositor is running — never on the visual critical path
After=playos-compositor.service dbus.service
Wants=dbus.service
PartOf=playos-async.target

[Service]
Type=dbus
BusName=org.bluez
ExecStart=/usr/lib/bluetooth/bluetoothd --noplugin=sap
Restart=on-failure
RestartSec=3
Nice=5
CPUWeight=500
IOWeight=500

[Install]
WantedBy=playos-async.target
