# archiso/profiles/playos/airootfs/etc/systemd/system/playos-async.target
[Unit]
Description=PlayOS Async Background Services
Documentation=https://github.com/PlayOS-Foundation/playos-refdistro
# Async services start independently; do NOT wait for the compositor.
# If the compositor crashes, networking + sshd must still come up
# so the device can be debugged remotely.
After=network.target
Wants=playos-audio.service
Wants=playos-network.service
Wants=playos-bluetooth.service
Wants=playos-library.service
Wants=playos-update.service
