# archiso/profiles/playos/airootfs/etc/systemd/system/playos-firstboot.service
[Unit]
Description=PlayOS system bootstrap (network, SSH, locale)
Before=systemd-firstboot.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/timedatectl set-timezone UTC
ExecStart=/usr/bin/localectl set-locale LANG=en_US.UTF-8
ExecStart=/usr/bin/localectl set-keymap us
ExecStart=/usr/bin/systemctl enable --now NetworkManager
ExecStart=/usr/bin/systemctl enable --now sshd

[Install]
WantedBy=sysinit.target
